# vdr/debian/vdr.postrm
#! /bin/sh

case "$1" in
    purge|remove)
        update-rc.d -f vdr remove || true
        
        if [ "$1" = 'purge' ]
        then
            rm -rf /etc/vdr
            rm -rf /var/cache/vdr
            rm -rf /var/lib/video.00
            rm -f /var/run/vdr.pid
        fi
    ;;
esac
