# vdr/debian/vdr.prerm
#!/bin/sh

case "$1" in
    remove)
        if [ -x "/etc/init.d/vdr" ]
        then
            invoke-rc.d vdr stop || true
        fi
    ;;
esac

